# target/linux/ath79/generic/base-files/etc/uci-defaults/10_wifi_setting
. /lib/functions.sh

case "$(board_name)" in	
    npotelecom,cpe-wifi-2g)
	uci set wireless.@wifi-device[0].disabled="0"
        uci set wireless.@wifi-iface[0].disabled="0"
        uci set wireless.@wifi-iface[0].ssid="npotelecom-wifi-2g"
        uci set wireless.@wifi-iface[0].key="password"
        uci set wireless.@wifi-iface[0].encryption="psk2"
        uci set wireless.@wifi-device[0].country="RU"
        uci set wireless.@wifi-device[0].channel="1"
        uci commit wireless
		;;
    npotelecom,cpe-wifi-5g)
		uci set wireless.@wifi-device[0].disabled="0"
        uci set wireless.@wifi-iface[0].disabled="0"
        uci set wireless.@wifi-iface[0].ssid="npotelecom-wifi-5g"
        uci set wireless.@wifi-iface[0].key="password"
        uci set wireless.@wifi-iface[0].encryption="psk2"
        uci set wireless.@wifi-device[0].country="RU"
        uci set wireless.@wifi-device[0].channel="36"
        uci commit wireless
		;;
esac

exit 0
