# target/linux/ramips/mt7621/base-files/etc/uci-defaults/99-network_hostname_default
uci -q batch << EOI

set network.lan.device='lan'
commit network

set wireless.@wifi-device[0].disabled='0'
set wireless.@wifi-iface[0].ssid='npotelecom-wifi-2g'
set wireless.@wifi-iface[0].cell_density='0'
set wireless.@wifi-iface[0].encryption='password'
set wireless.@wifi-iface[0].country='RU'

set wireless.@wifi-device[1].disabled='0'
set wireless.@wifi-iface[1].ssid='npotelecom-wifi-5g'
set wireless.@wifi-iface[1].cell_density='0'
set wireless.@wifi-iface[1].encryption='password'
set wireless.@wifi-iface[1].country='RU'

commit wireless

set system.@system[0].hostname='AP5000AX'
commit system

EOI
